# .config/autostart/albert.desktop
[Desktop Entry]
Type=Application
Exec=albert
Hidden=false
X-MATE-Autostart-enabled=true
Name[pt_BR]=Albert
Name=Albert
Comment[pt_BR]=
Comment=
