# .config/autostart/clipit.desktop
[Desktop Entry]
Type=Application
Exec=clipit
Hidden=false
X-MATE-Autostart-enabled=true
Name[pt_BR]=ClipIt
Name=ClipIt
Comment[pt_BR]=
Comment=
