# home/lucas/.config/systemd/user/mega.timer
[Unit]
Description=Mega Sync Timer (runs every hour)

[Timer]
OnUnitActiveSec=1h
Persistent=true
Unit=mega.service

[Install]
WantedBy=timers.target
